# cron_tests/crontab

* * * * * /bin/bash /root/run_tests.sh

