# app/crontab

* * * * * /bin/bash /var/www/run_tests.sh

